# deploy/dashboard-startup.service
[Unit]
Description=Dashboard gastos startup updater
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/pi/dashboard-gastos/scripts/run_dashboard_pi.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
